# smf/manifests/haproxy.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type="manifest" name="haproxy">
    <service name="manta/haproxy" type="service" version="1">

	<create_default_instance enabled="false" />
	<single_instance />

	<dependency name="network"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/milestone/network:default" />
	</dependency>

	<dependency name="filesystem"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/system/filesystem/local" />
	</dependency>

        <dependency name="mdata"
                    grouping="require_all"
                    restart_on="none"
                    type="service">
            <service_fmri value="svc:/smartdc/mdata:execute" />
        </dependency>

	<exec_method type="method"
		     name="start"
		     exec="/opt/local/sbin/haproxy -f /opt/smartdc/moray/etc/haproxy.cfg -D"
		     timeout_seconds="30">
	    <method_context working_directory="/opt/smartdc/moray">
		<method_environment>
		    <envvar name="PATH"
			    value="/opt/local/bin:/usr/bin:/usr/sbin:/bin"/>
		    <envvar name="LD_PRELOAD_32"
			    value="/usr/lib/extendedFILE.so.1" />
		</method_environment>
	    </method_context>
	</exec_method>

	<exec_method type="method"
		     name="stop"
		     exec=":kill"
		     timeout_seconds="30" />

	<template>
	    <common_name>
		<loctext xml:lang="C">HAProxy</loctext>
	    </common_name>
	</template>

    </service>
</service_bundle>
